# docs/systemd/rnaa-run.timer
[Unit]
Description=Periodic RNAA orchestration

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Unit=rnaa-run.service

[Install]
WantedBy=timers.target
